# systemd/user/almanac-hermes-docs-sync.timer
[Unit]
Description=Refresh synced Hermes documentation in the Almanac vault every 4 hours

[Timer]
OnBootSec=10m
OnUnitActiveSec=4h
Unit=almanac-hermes-docs-sync.service

[Install]
WantedBy=timers.target
