# systemd/user/almanac-curator-discord-onboarding.service
[Unit]
Description=Almanac Curator Discord onboarding worker
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/almanac/bin/almanac-curator-discord-onboarding.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
