# systemd/user/almanac-memory-synth.timer
[Unit]
Description=Synthesize compact Almanac memory cards every 30 minutes

[Timer]
OnBootSec=6m
OnUnitActiveSec=30m
Unit=almanac-memory-synth.service

[Install]
WantedBy=timers.target
