# systemd/user/almanac-ssot-batcher.timer
[Unit]
Description=Process queued Almanac SSOT webhook events every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=1m
Unit=almanac-ssot-batcher.service

[Install]
WantedBy=timers.target
