# systemd/user/almanac-notification-delivery.service
[Unit]
Description=Almanac notification delivery worker

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/almanac-notification-delivery.sh
