# systemd/user/almanac-qmd-update.service
[Unit]
Description=Almanac qmd index refresh

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/qmd-refresh.sh
