# systemd/user/almanac-github-backup.service
[Unit]
Description=Almanac vault backup to GitHub

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/backup-to-github.sh
