# systemd/user/almanac-quarto-render.service
[Unit]
Description=Almanac Quarto render

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/render-quarto.sh
