# systemd/user/almanac-pdf-ingest-watch.service
[Unit]
Description=Compatibility wrapper for the Almanac vault watcher
After=network.target

[Service]
ExecStart=%h/almanac/bin/pdf-ingest-watch.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
