# systemd/user/almanac-curator-refresh.timer
[Unit]
Description=Run the Almanac Curator refresh every hour

[Timer]
OnBootSec=2m
OnUnitActiveSec=1h
Unit=almanac-curator-refresh.service

[Install]
WantedBy=timers.target
