# systemd/user/almanac-memory-synth.service
[Unit]
Description=Almanac semantic memory synthesis

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/memory-synth.sh
