# systemd/user/almanac-pdf-ingest.service
[Unit]
Description=Almanac PDF ingestion and qmd refresh

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/pdf-ingest.sh
