# systemd/user/almanac-mcp.service
[Unit]
Description=Almanac control-plane MCP server
After=network.target

[Service]
ExecStart=%h/almanac/bin/almanac-mcp-server.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
