# systemd/user/almanac-curator-refresh.service
[Unit]
Description=Almanac Curator refresh job

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/curator-refresh.sh
