# systemd/user/almanac-hermes-docs-sync.service
[Unit]
Description=Sync Hermes documentation into the Almanac vault

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/sync-hermes-docs-into-vault.sh
