# systemd/user/almanac-ssot-batcher.service
[Unit]
Description=Almanac SSOT / Notion batch processor

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/almanac-ssot-batcher.sh
