# systemd/user/almanac-health-watch.service
[Unit]
Description=Almanac health watch operator notifier

[Service]
Type=oneshot
ExecStart=%h/almanac/bin/health-watch.sh
