# systemd/user/arclink-ssot-batcher.service
[Unit]
Description=ArcLink SSOT / Notion batch processor

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/arclink-ssot-batcher.sh
