# systemd/user/arclink-health-watch.service
[Unit]
Description=ArcLink health watch operator notifier

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/health-watch.sh
