# systemd/user/arclink-qmd-update.service
[Unit]
Description=ArcLink qmd index refresh

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/qmd-refresh.sh
