# systemd/user/arclink-notification-delivery.timer
[Unit]
Description=Deliver queued ArcLink notifications promptly

[Timer]
OnBootSec=5s
OnUnitActiveSec=5s
AccuracySec=1s
Unit=arclink-notification-delivery.service

[Install]
WantedBy=timers.target
