# systemd/user/arclink-curator-refresh.service
[Unit]
Description=ArcLink Curator refresh job

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/curator-refresh.sh
