# systemd/user/arclink-quarto-render.service
[Unit]
Description=ArcLink Quarto render

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/render-quarto.sh
