# systemd/user/arclink-memory-synth.service
[Unit]
Description=ArcLink semantic memory synthesis

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/memory-synth.sh
