# systemd/user/arclink-notification-delivery.timer
[Unit]
Description=Deliver queued ArcLink notifications every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=1m
Unit=arclink-notification-delivery.service

[Install]
WantedBy=timers.target
