# systemd/user/arclink-notification-delivery.service
[Unit]
Description=ArcLink notification delivery worker

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/arclink-notification-delivery.sh
