# systemd/user/arclink-curator-refresh.timer
[Unit]
Description=Run the ArcLink Curator refresh every hour

[Timer]
OnBootSec=2m
OnUnitActiveSec=1h
Unit=arclink-curator-refresh.service

[Install]
WantedBy=timers.target
