# systemd/user/arclink-pdf-ingest.service
[Unit]
Description=ArcLink PDF ingestion and qmd refresh

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/pdf-ingest.sh
