# systemd/user/arclink-github-backup.service
[Unit]
Description=ArcLink vault backup to GitHub

[Service]
Type=oneshot
ExecStart=%h/arclink/bin/backup-to-github.sh
